# systemd/mateusz-apps-certbot-renew.timer
[Unit]
Description=Check the Mateusz Apps public TLS certificate twice daily

[Timer]
OnCalendar=*-*-* 00,12:17:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
